# dotfiles/.crontab
0 0 * * * sudo apt update && sudo apt -y upgrade && sudo apt -y autoremove && rm -rf /home/aaron/.local/share/Trash/*
